# meta-sw/recipes-sw/errorApp/files/src/error-app.service
[Unit]
Description=erro-app server
[Service]
Type=simple
ExecStart=/usr/bin/erro-app

[Install]
WantedBy=multi-user.target